# crontab
* * * * * ubuntu bash /home/ubuntu/apps/app-template/repository/cron.sh python update_feed.py 
